# osd/fa_osd_input_4.service
Description=FA OSD Input 4
Wants=network-online.target systemd-networkd-wait-online.service redis.service rabbitmq.service

[Service]
Restart=on-failure
RestartSec=3s
Environment="CONFIG_FILE=/opt/fcs/osd/input/input_4/.env"
ExecStart=/opt/fcs/osd/input/input_4/fa-fcs-osd

[Install]
WantedBy=multi-user.target

